# DW_ingest_radio.service
[Unit]
Description=Ingest Radion Comedy
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/DW_ingest_radio.sh
User=root
