# usenet_org.service
[Unit]
Description=usenet file organiser
After=network.target

[Service]
# Dynamically determine user based on your HOST logic
# We use 'dan' as a fallback, but systemd usually runs as a specific User.
# If you want it to run as root to allow the auto-apt install:
User=root
Group=root

Type=simple
# Adjust the path below to where you actually saved the .sh file
ExecStart=/bin/bash /usr/local/bin/usenet_org.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
